# crontab
0 6 * * * . /a/env && python /a/update_ranking.py >> /var/log/cron.log 2>&1
