# crontab
0 6 * * * /srv/rankingbot/update_ranking.sh >> /var/log/cron.log 2>&1
