# crontab
0 6 * * * . /a/env && /usr/local/bin/python /a/update_ranking.py >> /tmp/log 2>&1
